# systemd/container-cleanup.service
[Unit]
Description=Container image cleanup service

[Service]
Type=oneshot
ExecStart=%h/.containers/scripts/cleanup.sh
